# deploy/templates/rclone-archive.timer.j2
[Unit]
Description=Perform rclone archive

[Timer]
OnCalendar=02:30:00

[Install]
WantedBy=timers.target
